# docker/resources/supervisor/conf.d/crond.conf
[program:crond]
command=crond -f -c /etc/crontabs
stopsignal=KILL
autostart=true
autorestart=true
