# docker/resources/supervisor/conf.d/nginx-prom-exporter.conf
[program:nginx-prom-exporter]
command=nginx-prometheus-exporter --nginx.scrape-uri=http://127.0.0.1:80/nginx_status
stopsignal=KILL
stopasgroup=true
killasgroup=true
autostart=true
autorestart=true
startretries=5
numprocs=1
startsecs=0
stderr_logfile=/var/log/nginx/nginx-prometheus-exporter.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/var/log/nginx/nginx-prometheus-exporter.log
stdout_logfile_maxbytes=10MB
redirect_stderr=true
